# automation/sitemap-gen.path
[Unit]
Description=Watch /var/www/html for changes

[Path]
PathChanged=/var/www/html/
DelaySec=15

[Install]
WantedBy=multi-user.target
